# rig_backup/etc_systemd/ai-rig-weekly-restart.timer
[Unit]
Description=Weekly llama-server restart (Sunday 4 AM)

[Timer]
OnCalendar=Sun *-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target
