# rig_backup/etc_systemd/ai-rig-resume.service
[Unit]
Description=AI Rig - Post-Resume GPU Health Check + Auto-Reload
After=suspend.target

[Service]
Type=oneshot
ExecStart=/opt/ai-rig-resume.sh

[Install]
WantedBy=suspend.target
