# rig_backup/etc_systemd/ai-rig-watchdog.service
[Unit]
Description=AI Rig Watchdog (GPU health + DB sync + OOM)

[Service]
Type=oneshot
ExecStart=/opt/ai-rig-watchdog.sh
