# reingest_spec/systemd/cbic-rag-api.service.d/power_dpm.conf
# Force high power on all embed-pool + reranker cards BEFORE the API loads its pool.
#
# Why this exists (codified 2026-05-08, INCIDENTS_ARCHIVE.md):
#   The base systemd ExecStartPre on bge-reranker / qwen3 only forces card0 + card2 to
#   `dpm=high`. The embed pool cards {1,3,4,5,6} get forced to high during ingest_v2's
#   pool warmup, but revert to `dpm=auto` (sclk=800Mhz, low power) when ingest exits.
#   Subsequent API /retrieve queries embed via the pool at LOW power → 10-15s/query +
#   reranker timeouts at 30s under any concurrency.
#
# This drop-in re-asserts `high` on every card the API will touch (0,1,2,3,4,5,6) at
# every API start. It runs as root (systemd) so the writes succeed.
#
# Mirrored to: D:/_gpu_rig_ai/reingest_spec/systemd/cbic-rag-api.service.d/power_dpm.conf
#
# To install on the rig:
#   sudo install -m 0644 /tmp/cbic_api_power_dpm.conf /etc/systemd/system/cbic-rag-api.service.d/power_dpm.conf
#   sudo systemctl daemon-reload
#   sudo systemctl restart cbic-rag-api

[Service]
ExecStartPre=/bin/bash -c 'for c in 0 1 2 3 4 5 6; do echo high > /sys/class/drm/card$c/device/power_dpm_force_performance_level 2>/dev/null || true; done'
