# rig_backup/etc_modprobe/99-wol.rules
ACTION=="add", SUBSYSTEM=="net", NAME=="eth*", RUN+="/usr/sbin/ethtool -s %k wol g"
