# reingest_spec/systemd/bge-reranker.service.d/power_dpm.conf
# Drop-in for /etc/systemd/system/bge-reranker.service.d/power_dpm.conf
# Force GPU 0 power_dpm to 'high' BEFORE bge-reranker loads model.
# Same rationale as qwen3-14b.service.d/power_dpm.conf — codified 2026-05-07.
[Service]
ExecStartPre=/bin/bash -c "echo high > /sys/class/drm/card0/device/power_dpm_force_performance_level && sleep 2"
