# reingest_spec/systemd/qwen3-14b.service.d/power_dpm.conf
# Drop-in for /etc/systemd/system/qwen3-14b.service.d/power_dpm.conf
# Force GPU 2 power_dpm to 'high' BEFORE llama-server loads model.
#
# WHY: 2026-05-07 root-caused a 300x slowdown after reboot. qwen3-14b service
# auto-starts at boot via WantedBy=multi-user.target. If power_dpm is still at
# default (auto/low), the model loads while GPU 2 is in low power state, and
# never escapes — inference becomes ~104s per 10 tokens vs the expected ~330ms.
# Direct in-process llama-cpp-python had the SAME slowdown when GPU was in
# low-power, proving it was the GPU power state, not the HTTP/proxy path.
#
# WITHOUT this drop-in: every reboot reproduces the slowdown until manual
# intervention (echo high > .../power_dpm_force_performance_level).
# WITH this drop-in: power_dpm is forced high before llama-server invokes
# the GPU at all, so model loads at full clocks (2855 MHz on RX 6700 XT).
#
# Codified 2026-05-07 after RG-driven verification: direct in-process bench
# showed 27.4 tok/s, HTTP llama-server bench showed 36.14 tok/s — both fast
# only when power_dpm was set high BEFORE model load.
[Service]
ExecStartPre=/bin/bash -c "echo high > /sys/class/drm/card2/device/power_dpm_force_performance_level && sleep 2"
